# docker/customer/sources/supervisor/consul.conf
[program:consul]
command=/usr/local/bin/consul agent -config-dir /etc/consul.d
directory=/var/consul
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_logfile = /var/log/supervisor/%(program_name)s.err.log
priority=6
