# docker/customer/sources/supervisor/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile = /var/log/%(program_name)s.log
priority=5
