# docker/mysql/sources/supervisor/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld --init-file=/tmp/init.sql
redirect_stderr=true
stdout_syslog=true
autostart=true
autorestart=true
