# docker/base/sources/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr=true
stdout_syslog=true
autostart=true
autorestart=true
user=root
