# docker/log/sources/nginx-supervisor.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log