# vagrant/modules/zipkin_collector/files/web.conf
[program:web]
command=java -cp libs -jar zipkin-web-1.1.0.jar -f config/web-dev.scala -D local_docroot=/opt/app/zipkin-web/resources
directory=/opt/app/zipkin-web
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/zipkin/web.out.log
stderr_logfile=/var/log/zipkin/web.err.log
