# docker/log/sources/supervisor/logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash agent -f /etc/logstash/conf.d/
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
